# src/stats/config/crontab
# see https://crontab.guru/ for help on crontabs syntax

#   m     h     dom     mon     dow     command
   */5    *      *       *       *      /srv/src/config/run_stats.sh 
